# coprocessor/led/led.service
[Unit]
Description=Robot LED Control Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/photon/led
ExecStart=/home/photon/led/.venv/bin/python led_service.py
Restart=always
RestartSec=5
StartLimitIntervalSec=0

[Install]
WantedBy=multi-user.target  